# systemd/kvmind-register.timer
[Unit]
Description=KVMind Registration Timer

[Timer]
OnBootSec=30
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
