# systemd/kvmind.service
[Unit]
Description=KVMind Bridge Service
Documentation=https://github.com/kvmind-labs/kvmind-integration
After=network.target kvmd.service
Wants=kvmd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/kvmind/kdkvm
ExecStartPre=/bin/bash -c 'mount -o remount,rw /var/lib/kvmd/msd 2>/dev/null; mkdir -p /var/lib/kvmd/msd/.kdkvm 2>/dev/null; mount -o remount,ro /var/lib/kvmd/msd 2>/dev/null; true'
ExecStartPre=/bin/mkdir -p /var/log/kdkvm /var/log/nginx
ExecStart=/opt/kvmind/kdkvm/venv/bin/python -m lib.server
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kvmind

Environment=KVMIND_KVM_CONFIG=/etc/kdkvm/config.yaml
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/kvmind/kdkvm
EnvironmentFile=-/etc/kdkvm/ai.env

MemoryMax=256M
MemoryHigh=200M
CPUQuota=80%
TasksMax=64
ProtectHome=no

[Install]
WantedBy=multi-user.target
